# contrib/systemd/heimdall.service.d/01-authoritative.conf
# Heimdall systemd drop-in — authoritative role
#
# Deploy this drop-in alongside the base heimdall.service to enable the
# authoritative server role. Per-role drop-ins are stackable: a hybrid
# (authoritative + recursive + forwarder) instance applies all three.
#
# This drop-in mounts:
#   /etc/heimdall/           RO  — configuration
#   /etc/heimdall/tls/       RO  — TLS certificates (if DoT/DoH/DoQ enabled)
#   /var/lib/heimdall/zones/ RO  — authoritative zone files
#   /run/heimdall/           RW  — runtime state (already RW via RuntimeDirectory)

[Service]
ReadOnlyPaths=/etc/heimdall /var/lib/heimdall/zones
