# contrib/systemd/heimdall.service.d/02-recursive.conf
# Heimdall systemd drop-in — recursive resolver role
#
# Deploy this drop-in alongside the base heimdall.service to enable the
# recursive resolver role. RFC 5011 trust-anchor rollover requires write
# access to the trust-anchor file.
#
# This drop-in mounts:
#   /etc/heimdall/                       RO  — configuration
#   /etc/heimdall/tls/                   RO  — TLS certificates (if DoT/DoH/DoQ enabled)
#   /var/lib/heimdall/trust-anchors.xml  RW  — RFC 5011 trust-anchor (rollover writes)

[Service]
ReadOnlyPaths=/etc/heimdall
ReadWritePaths=/var/lib/heimdall/trust-anchors.xml
