# contrib/systemd/heimdall.service.d/03-forwarder.conf
# Heimdall systemd drop-in — forwarder role
#
# Deploy this drop-in alongside the base heimdall.service to enable the
# forwarder role.
#
# This drop-in mounts:
#   /etc/heimdall/      RO  — configuration
#   /etc/heimdall/tls/  RO  — TLS certificates (if DoT/DoH/DoQ outbound)

[Service]
ReadOnlyPaths=/etc/heimdall
